# scripts/gitops/gitops-controller.timer
[Unit]
Description=GitOps Controller Timer - Poll GitHub every 2 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=2min
AccuracySec=10s

[Install]
WantedBy=timers.target
