# scripts/gitops/gitops-webhook.service
[Unit]
Description=GitOps Webhook Listener - Receive GitHub push events
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/gitops/gitops-webhook.py
WorkingDirectory=/opt/gitops
Restart=always
RestartSec=5
EnvironmentFile=/etc/gitops/config.env

# Note: ProtectSystem is intentionally omitted — the gitops controller
# needs broad access to run pct push/exec, git, and docker commands.

[Install]
WantedBy=multi-user.target
